# data/waypaperd.service
[Unit]
Description=Waypaper slideshow daemon
Documentation=https://anufrievroman.gitbook.io/waypaper/
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=WAYPAPERD_INTERVAL=1800
ExecStart=/usr/bin/env waypaperd ${WAYPAPERD_INTERVAL}
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillMode=control-group
NoNewPrivileges=true

[Install]
WantedBy=default.target
